# wallcron.cron
0 * * * * /home/iago/Projects/dotfiles/wallscript.sh >> /home/iago/Projects/dotfiles/wallmac.log
